# {{cookiecutter.repo_name}}/data/config/supervisor/supervisor-{{ cookiecutter.project_slug }}.ini
[program:{{ cookiecutter.project_slug }}-dev]
command=python3 run.py --profile=dev --debug=True --log=debug
directory=${CWD}/../../../
stdout_logfile=data/logs/{{ cookiecutter.project_slug }}-console.log
stderr_logfile=data/logs/{{ cookiecutter.project_slug }}-error.log
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=60
priority=999
